# doc/supervisord.conf
[program:mybo_es]
command=gunicorn -w4 -b0.0.0.0:8090 mybo_es.wsgi
directory=/usr/share/nginx/html/mybo_es
startsecs=0
stopwaitsecs=0
autostart=true
autorestart=true